# overlay/udev/99-snd.rules
SUBSYSTEM=="sound", ACTION=="add", RUN+="/bin/systemctl restart snapclient.service"
